# extras/avalanchemq-amqp.socket
[Unit]
Description=AvalancheMQ

[Socket]
ListenStream=0.0.0.0:15672
Service=avalanchemq.service

[Install] 
WantedBy=sockets.target
