# extras/avalanchemq-http.socket
[Unit]
Description=AvalancheMQ

[Socket]
ListenStream=/tmp/avalanchemq.sock

Service=avalanchemq.service

[Install] 
WantedBy=sockets.target
