# extras/lavinmq-http.socket
[Unit]
Description=LavinMQ

[Socket]
ListenStream=/tmp/lavinmq.sock

Service=lavinmq.service

[Install]
WantedBy=sockets.target
