# extras/lavinmq-amqp.socket
[Unit]
Description=LavinMQ

[Socket]
ListenStream=0.0.0.0:15672
Service=lavinmq.service

[Install]
WantedBy=sockets.target
